# io.arslan.dark-mode-notify.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
"http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.arslan.dark-mode-notify</string>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/fatih/.dark-mode-notify-stderr.log</string>
    <key>StandardOutPath</key>
    <string>/Users/fatih/.dark-mode-notify-stdout.log</string>
    <key>ProgramArguments</key>
    <array>
       <string>/usr/local/bin/dark-mode-notify</string>
       <string>/usr/local/bin/fish</string>
       <string>-c</string>
       <string>change_background</string>
    </array>
</dict>
</plist>
